# deploy/onlybots-drift-check.service
[Unit]
Description=OnlyBots Drift Check (Phase 4) — queues re-runs for stale verified services
After=network.target postgresql.service
# Only run this when the main verifier is healthy — no point queuing runs if
# nothing is polling.
Requires=onlybots-verifier.service

[Service]
Type=oneshot
User=onlybots
Group=onlybots
WorkingDirectory=/opt/onlybots/verifier
EnvironmentFile=/opt/onlybots/.env
EnvironmentFile=/opt/onlybots/verifier/.env
# --drift-check queues re-runs; the main verifier daemon picks them up from
# the normal polling loop. This service exits as soon as queuing is done.
ExecStart=/opt/onlybots/verifier/venv/bin/python main.py --drift-check --staleness-hours=24

# Same hardening as the main verifier — drift-check doesn't run agent
# code itself, but it's belt-and-suspenders to match the security posture.
NoNewPrivileges=true
PrivateTmp=true
IPAddressAllow=127.0.0.0/8
IPAddressDeny=169.254.0.0/16 10.0.0.0/8 172.16.0.0/12 192.168.0.0/16 100.64.0.0/10 fc00::/7 fe80::/10
